# docker/ckan/crontab
*/15 * * * *    /usr/lib/ckan/venv/bin/ckan -c /etc/ckan/production.ini harvester run
